# .config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.3.24)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/deder/.nvm/versions/node/v24.14.1/bin/node /home/deder/.nvm/versions/node/v24.14.1/lib/node_modules/openclaw/dist/index.js gateway --port 18789
ExecStartPost=/bin/bash -c 'sleep 60 && /home/deder/.nvm/versions/node/v24.14.1/bin/node /home/deder/.nvm/versions/node/v24.14.1/lib/node_modules/openclaw/dist/index.js cron run --name brewbird-daily-refresh &'
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME=/home/deder
Environment=TMPDIR=/tmp
Environment=NODE_EXTRA_CA_CERTS=/etc/ssl/certs/ca-certificates.crt
Environment=PATH=/home/deder/.nvm/versions/node/v24.14.1/bin:/home/deder/.nvm/current/bin:/home/deder/.local/bin:/home/deder/.npm-global/bin:/home/deder/bin:/home/deder/.volta/bin:/home/deder/.asdf/shims:/home/deder/.bun/bin:/home/deder/.fnm/current/bin:/home/deder/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.24

[Install]
WantedBy=default.target
